# example/watch-ban.upstart
description "Watch domain lists to generate DNS black-hole configurations"
author  "John Manero <john.manero@gmail.com>"

start on startup
stop on shutdown

respawn
respawn limit 5 10

console log
exec /opt/watch-ban/bin/server -c /opt/watch-ban/example/config.json
